# deploy/free-oracle-tailscale/grocery-cockpit.service
[Unit]
Description=Grocery Cockpit
After=network-online.target grocery-display.service
Wants=network-online.target grocery-display.service

[Service]
Type=simple
WorkingDirectory=/opt/grocery-cockpit
Environment=DISPLAY=:99
Environment=GROCERY_CHROME_PATH=/usr/bin/chromium
ExecStart=/usr/bin/python3 -u /opt/grocery-cockpit/grocery_cockpit.py serve --host 127.0.0.1
Restart=always
RestartSec=5
User=grocery
Group=grocery

[Install]
WantedBy=multi-user.target
