# deploy/free-oracle-tailscale/grocery-x11vnc.service
[Unit]
Description=Grocery Cockpit VNC bridge
After=grocery-display.service
Requires=grocery-display.service

[Service]
Type=simple
ExecStart=/usr/bin/x11vnc -display :99 -localhost -forever -shared -nopw -rfbport 5900
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
